# tradingbot.service
[Unit]
Description=AI Trade Planner Telegram Bot
After=network.target

[Service]
Type=simple
# Replace with your actual VPS username
User=your_vps_username
WorkingDirectory=/home/your_vps_username/trade-bot
EnvironmentFile=/home/your_vps_username/trade-bot/.env
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
# Capture stdout/stderr to journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
